# upstart/lager-server.conf
description "Lager Web Server"

start on runlevel [2345]
stop on runlevel [^2345]

respawn

kill timeout 30

script
/bin/bash <<EOT
    /home/pi/git/lager/bin/lager-server
EOT
end script
